# metabase/rootfs/etc/s6-overlay/s6-rc.d/metabase/run
#!/usr/bin/with-contenv bashio
set -e

JAVA_MEM=$(bashio::config 'java_memory')

bashio::log.info "Starting Metabase (heap: ${JAVA_MEM})..."

# Metabase internal DB (H2) persists in /data so config survives restarts
export MB_DB_FILE=/data/metabase.db
export MB_DB_TYPE=h2
export MB_JETTY_HOST=127.0.0.1
export MB_JETTY_PORT=3000

# Disable Metabase update checks and analytics
export MB_CHECK_FOR_UPDATES=false
export MB_ANON_TRACKING_ENABLED=false

# Disable password complexity requirements (admin password is auto-managed)
export MB_PASSWORD_COMPLEXITY=weak
export MB_PASSWORD_LENGTH=8

exec java \
    -Xms256m \
    -Xmx"${JAVA_MEM}" \
    -jar /opt/metabase/metabase.jar
