# oxigraph/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Generating nginx configuration..."

AUTH_TYPE=$(bashio::config 'auth_type')
NETWORK_ACCESS=$(bashio::config 'network_access')

# Remove any config from a previous run
rm -f /etc/nginx/http.d/oxigraph*.conf

# ── Ingress server (HA supervisor only, port 8099, no auth) ──────────────────
# ingress_entry is "" so HA supervisor forwards GET / correctly (no double slash)
cat > /etc/nginx/http.d/oxigraph-ingress.conf << 'NGINX'
server {
    listen 8099 default_server;

    # Only allow the HA supervisor proxy; reject all other sources
    allow  172.30.32.2;
    deny   all;

    location / {
        proxy_pass         http://127.0.0.1:7879/;
        proxy_http_version 1.1;
        proxy_set_header   Host            $http_host;
        proxy_set_header   Upgrade         $http_upgrade;
        proxy_set_header   Connection      "upgrade";
        proxy_set_header   X-Real-IP       $remote_addr;
        proxy_set_header   X-Forwarded-For $proxy_add_x_forwarded_for;
        proxy_buffering    off;
    }
}
NGINX

# ── External server (port 7878, optional) ────────────────────────────────────
if bashio::config.true 'network_access'; then
    case "$AUTH_TYPE" in

        basic)
            USERNAME=$(bashio::config 'username')
            PASSWORD=$(bashio::config 'password')
            if [ -z "$USERNAME" ] || [ -z "$PASSWORD" ]; then
                bashio::log.fatal "auth_type=basic requires both 'username' and 'password' to be set!"
                exit 1
            fi
            htpasswd -bc /etc/nginx/.htpasswd "$USERNAME" "$PASSWORD"
            bashio::log.info "External access enabled – Basic auth (user: ${USERNAME})"
            cat > /etc/nginx/http.d/oxigraph-external.conf << 'NGINX'
server {
    listen 7878 default_server;
    location / {
        auth_basic           "Oxigraph";
        auth_basic_user_file /etc/nginx/.htpasswd;
        proxy_pass         http://127.0.0.1:7879/;
        proxy_http_version 1.1;
        proxy_set_header   Host            $http_host;
        proxy_set_header   Upgrade         $http_upgrade;
        proxy_set_header   Connection      "upgrade";
        proxy_set_header   X-Real-IP       $remote_addr;
        proxy_set_header   X-Forwarded-For $proxy_add_x_forwarded_for;
        proxy_buffering    off;
    }
}
NGINX
            ;;

        bearer)
            BEARER_TOKEN=$(bashio::config 'bearer_token')
            if [ -z "$BEARER_TOKEN" ]; then
                bashio::log.fatal "auth_type=bearer requires 'bearer_token' to be set!"
                exit 1
            fi
            bashio::log.info "External access enabled – Bearer token auth"
            # map directive must be in http context (http.d files are included there)
            printf 'map $http_authorization $oxigraph_ok {\n    "Bearer %s" 1;\n    default 0;\n}\n' \
                "$BEARER_TOKEN" > /etc/nginx/http.d/oxigraph-external.conf
            cat >> /etc/nginx/http.d/oxigraph-external.conf << 'NGINX'
server {
    listen 7878 default_server;
    location / {
        # Return 401 when the Bearer token is missing or wrong.
        # proxy_pass is in the outer location and executes only when the if-block does NOT match.
        if ($oxigraph_ok = 0) {
            add_header WWW-Authenticate 'Bearer realm="Oxigraph"' always;
            return 401;
        }
        proxy_pass         http://127.0.0.1:7879/;
        proxy_http_version 1.1;
        proxy_set_header   Host            $http_host;
        proxy_set_header   Upgrade         $http_upgrade;
        proxy_set_header   Connection      "upgrade";
        proxy_set_header   X-Real-IP       $remote_addr;
        proxy_set_header   X-Forwarded-For $proxy_add_x_forwarded_for;
        proxy_buffering    off;
    }
}
NGINX
            ;;

        none|*)
            bashio::log.warning "⚠️  External access is enabled WITHOUT authentication!"
            bashio::log.warning "Anyone who can reach port 7878 can read AND modify your triple store."
            bashio::log.warning "Set auth_type to 'basic' or 'bearer' to protect Oxigraph."
            cat > /etc/nginx/http.d/oxigraph-external.conf << 'NGINX'
server {
    listen 7878 default_server;
    location / {
        proxy_pass         http://127.0.0.1:7879/;
        proxy_http_version 1.1;
        proxy_set_header   Host            $http_host;
        proxy_set_header   Upgrade         $http_upgrade;
        proxy_set_header   Connection      "upgrade";
        proxy_set_header   X-Real-IP       $remote_addr;
        proxy_set_header   X-Forwarded-For $proxy_add_x_forwarded_for;
        proxy_buffering    off;
    }
}
NGINX
            ;;
    esac
else
    bashio::log.info "External network access is disabled (network_access=false). Only HA ingress is active."
fi

bashio::log.info "Starting nginx..."
exec nginx -g "daemon off;"
