# prefect-server/rootfs/etc/s6-overlay/s6-rc.d/prefect-server/run
#!/command/with-contenv bashio
# ==============================================================================
# Start Prefect Server
# ==============================================================================

bashio::log.info "Preparing Prefect Server data directory..."
mkdir -p /data/prefect

# Prefect's block auto-registration writes to PREFECT_HOME.
# Ensure the directory exists before the server starts.
mkdir -p /root/.prefect 2>/dev/null || true

bashio::log.info "Starting Prefect Server on 127.0.0.1:4201..."
exec /opt/prefect/bin/prefect server start \
    --host 127.0.0.1 \
    --port 4201
