# metabase/rootfs/etc/s6-overlay/s6-rc.d/db-sync/run
#!/usr/bin/with-contenv bashio
set -e

HA_SQLITE="/config/home-assistant_v2.db"
DB_COPY="/data/ha_recorder.db"
INTERVAL=300  # 5 minutes

if [ ! -f "$HA_SQLITE" ]; then
    bashio::log.info "db-sync: no SQLite DB at ${HA_SQLITE}, sleeping."
    exec sleep infinity
fi

bashio::log.info "db-sync: refreshing ${DB_COPY} from ${HA_SQLITE} every ${INTERVAL}s."

while true; do
    sleep "$INTERVAL"
    sqlite3 "$HA_SQLITE" ".backup '${DB_COPY}'"
    bashio::log.debug "db-sync: snapshot refreshed."
done
