# metabase/rootfs/etc/s6-overlay/s6-rc.d/db-sync/run
#!/usr/bin/with-contenv bashio
set -e

HA_SQLITE="/config/home-assistant_v2.db"
DB_COPY="/data/ha_recorder.db"
INTERVAL=300  # 5 minutes

if [ ! -f "$HA_SQLITE" ]; then
    bashio::log.info "db-sync: no SQLite DB at ${HA_SQLITE}, sleeping."
    exec sleep infinity
fi

bashio::log.info "db-sync: refreshing ${DB_COPY} from ${HA_SQLITE} every ${INTERVAL}s."

while true; do
    sleep "$INTERVAL"
    # VACUUM INTO only needs read access; no lock files on the ro mount.
    # Write to a temp file first, then atomically replace the copy so
    # Metabase never sees a missing or half-written database file.
    TEMP="${DB_COPY}.tmp"
    rm -f "$TEMP"
    if sqlite3 "file:${HA_SQLITE}?mode=ro&immutable=1" "VACUUM INTO '${TEMP}';"; then
        mv -f "$TEMP" "$DB_COPY"
        bashio::log.debug "db-sync: snapshot refreshed."
    else
        bashio::log.warning "db-sync: VACUUM INTO failed, keeping previous copy."
        rm -f "$TEMP"
    fi
done
