# oxigraph/rootfs/etc/s6-overlay/s6-rc.d/oxigraph/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting Oxigraph SPARQL server on 127.0.0.1:7879..."
mkdir -p /data/storage

exec /usr/local/bin/oxigraph serve \
    --location /data/storage \
    --bind 127.0.0.1:7879
