# metabase/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio
set -e

MB_URL="http://127.0.0.1:3000"

# ── Wait for Metabase to become healthy ──────────────────────────────────────
bashio::log.info "Waiting for Metabase to start..."
for i in $(seq 1 120); do
    if curl -sf "${MB_URL}/api/health" | jq -e '.status == "ok"' >/dev/null 2>&1; then
        bashio::log.info "Metabase is healthy."
        break
    fi
    if [ "$i" -eq 120 ]; then
        bashio::log.fatal "Metabase did not become healthy within 10 minutes."
        exit 1
    fi
    sleep 5
done

# ── Bootstrap (idempotent) ───────────────────────────────────────────────────
/usr/bin/bootstrap.sh

# ── Obtain session token for SSO ─────────────────────────────────────────────
ADMIN_EMAIL=$(jq -r .email /data/admin.json)
ADMIN_PASS=$(jq -r .password /data/admin.json)

get_session_token() {
    curl -sf -X POST "${MB_URL}/api/session" \
        -H "Content-Type: application/json" \
        -d "{\"username\":\"${ADMIN_EMAIL}\",\"password\":\"${ADMIN_PASS}\"}" \
    | jq -r .id
}

SESSION_TOKEN=$(get_session_token)
if [ -z "$SESSION_TOKEN" ] || [ "$SESSION_TOKEN" = "null" ]; then
    bashio::log.fatal "Could not obtain Metabase session token."
    exit 1
fi
bashio::log.info "Metabase session obtained."

# ── Ensure theme CSS file exists (theme-sync service keeps it updated) ───────
if [ ! -f /data/ha-theme.css ]; then
    echo "" > /data/ha-theme.css
fi

# ── Write nginx config ──────────────────────────────────────────────────────
rm -f /etc/nginx/conf.d/*.conf

cat > /etc/nginx/conf.d/metabase-ingress.conf << NGINX
server {
    listen 8099 default_server;

    # Only allow the HA supervisor proxy
    allow  172.30.32.2;
    deny   all;

    client_max_body_size 64m;

    # Serve theme CSS override
    location = /ha-theme.css {
        alias /data/ha-theme.css;
        add_header Cache-Control "no-cache";
    }

    location / {
        proxy_pass         ${MB_URL}/;
        proxy_http_version 1.1;
        proxy_set_header   Host            \$http_host;
        proxy_set_header   Upgrade         \$http_upgrade;
        proxy_set_header   Connection      "upgrade";
        proxy_set_header   X-Real-IP       \$remote_addr;
        proxy_set_header   X-Forwarded-For \$proxy_add_x_forwarded_for;
        proxy_buffering    off;
        proxy_read_timeout 300s;

        # Inject session cookie for SSO – every ingress response sets the cookie
        add_header Set-Cookie "metabase.SESSION=${SESSION_TOKEN}; Path=/; HttpOnly; SameSite=Lax" always;

        # Inject theme CSS into Metabase HTML pages
        sub_filter '</head>' '<link rel="stylesheet" href="/ha-theme.css"></head>';
        sub_filter_once on;

        # Remove X-Frame-Options so HA can embed in its own frame
        proxy_hide_header X-Frame-Options;
    }
}
NGINX

bashio::log.info "nginx configuration generated."

# ── Session refresher (background) ──────────────────────────────────────────
# Refresh the session token every 6 hours to prevent expiry
(
    while true; do
        sleep 21600
        NEW_TOKEN=$(get_session_token)
        if [ -n "$NEW_TOKEN" ] && [ "$NEW_TOKEN" != "null" ]; then
            sed -i "s/metabase\.SESSION=[^;]*/metabase.SESSION=${NEW_TOKEN}/" \
                /etc/nginx/conf.d/metabase-ingress.conf
            nginx -s reload 2>/dev/null || true
            bashio::log.info "Metabase session token refreshed."
        fi
    done
) &

bashio::log.info "Starting nginx..."
exec nginx -g "daemon off;"
